# hermes-server/root/etc/systemd/system/oauth2-proxy.socket
[Unit]
Description=OAuth2 Proxy Socket

[Socket]
ListenStream=%t/oauth2.sock
SocketGroup=caddy
SocketMode=0660

[Install]
WantedBy=sockets.target
